# RC0592_ConfigurationManagement/configurationmanagement-demo/cmdemo2/cmdemo2daemon.service.in
[Unit]
Description=CM demo2 daemon

[Service]
Type=simple
ExecStart=@bindir@/cmdemo2daemon
ExecReload=/bin/kill -HUP $MAINPID
